# infra/deploy/simon-traps.service
[Unit]
Description=SIMON - Receptor de SNMP traps (UDP/162)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/monitoreo/monitor
ExecStart=/opt/monitoreo/monitor/.venv/bin/python /opt/monitoreo/monitor/trap_listener.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
